# recipes-dreambox/enigma2-plugins/enigma2-plugin-extensions-kodimediacenter/xbmc.service.in
[Unit]
Description=Start a script
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/bin/start-xbmc-standalone

[Install]
WantedBy=multi-user.target
